# mac-agent/Resources/com.lamp.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lamp.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/lamp-agent</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/lamp-agent.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/lamp-agent.log</string>
</dict>
</plist>
